# systemd/auraspeak.service
[Unit]
Description=AuraSpeak Server
After=network.target postgresql.service

[Service]
Type=simple
User=auraspeak
Group=auraspeak
WorkingDirectory=/opt/auraspeak/backend
Environment=NODE_ENV=production
ExecStart=/opt/auraspeak/backend/auraspeak
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 