# systemd/next-app.service
[Unit]
Description=Next.js Static App
After=network.target

[Service]
Type=simple
ExecStart=/opt/next-app/start.sh
Restart=always
User=www-data
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
